# docker/root/etc/s6-overlay/s6-rc.d/svc-heartrate/run
#!/usr/bin/with-contenv sh
exec s6-setuidgid abc gunicorn \
    --bind 0.0.0.0:5000 \
    --workers 2 \
    --chdir /app \
    webapp:app
